# rootfs/etc/services.d/hass-nec-control/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop the NEC TV Control service
# ==============================================================================

bashio::log.info "Stopping NEC TV Control service" 